# archcraft/etc/systemd/system/wallchange.timer
[Unit]
Description= Runs the wallchange.service job

[Timer]
OnUnitActiveSec=60s
OnBootSec=10s

[Install]
WantedBy=timers.target
